# celo-dashboard.service
[Unit]
Description=Celo Trader Dashboard (FastAPI + React)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/celo_trader
# Uvicorn serves the FastAPI backend.
# The built React frontend is served as static files by FastAPI itself.
ExecStart=/opt/celo_trader/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8501 --workers 1
Restart=always
RestartSec=5
StandardOutput=append:/opt/celo_trader/log/dashboard.log
StandardError=append:/opt/celo_trader/log/dashboard.log
Environment="PYTHONPATH=/opt/celo_trader"

[Install]
WantedBy=multi-user.target
